# data/export/upstart/master.conf.erb
PATH=/home/<%= user %>/.rbenv/shims:/home/<%= user %>/.rbenv/bin:$PATH
start on runlevel [2345]
stop on runlevel [!2345]
