# data/export/upstart/process_master.conf.erb
PATH=/home/<%= user %>/.rbenv/shims:/home/<%= user %>/.rbenv/bin:$PATH

start on starting <%= app %>
stop on stopping <%= app %>
